# deploy/sub2api-autodeploy.timer
[Unit]
Description=Optional polling fallback for Sub2API releases

[Timer]
OnBootSec=5min
OnUnitInactiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=sub2api-autodeploy.service

[Install]
WantedBy=timers.target
